# terraform/trash/libreddit/libreddit.service.j2
[Unit]
Description=libreddit Service
After=network.target

[Service]
ExecStart=/usr/local/bin/libreddit
Restart=always

[Install]
WantedBy=multi-user.target















